# docker/php/supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/log/supervisord.pid

[program:php-fpm]
command=docker-php-entrypoint php-fpm

[program:cron]
command=cron -f
autorestart=true

[program:horizon]
command=php /var/www/html/artisan horizon
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/horizon.log
